# profile.d/telegraf.sh
# Telegraf
export TELEGRAF_ENGENIUS_WIFI_COMMUNITY=public
export TELEGRAF_FORTINET_FORTIOS_COMMUNITY=public
export TELEGRAF_GENERIC_COMMUNITY=public
export TELEGRAF_HP_FUTURESMART_COMMUNITY=public
export TELEGRAF_JUNIPER_EX_COMMUNITY=public
export TELEGRAF_JUNIPER_SRX_COMMUNITY=public
export TELEGRAF_LENOVO_EMC_COMMUNITY=public
export TELEGRAF_UBIQUITI_EDGEOS_COMMUNITY=public
export TELEGRAF_VMWARE_ESXI_COMMUNITY=public
export TELEGRAF_VYOS_COMMUNITY=public
